# debian/popsec-daemon.service
[Unit]
Description=PopSec Daemon

[Service]
ExecStart=/usr/lib/popsec/popsec-daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
